# dns-nas-utils/DEBIAN/postrm
#!/bin/sh -e

if [ "$1" = "purge" ] ; then
    update-rc.d fan_daemon remove >/dev/null
fi
