# dns-nas-utils/DEBIAN/postinst
#!/bin/sh -e

if [ -x "/etc/init.d/fan-daemon" ]; then
    update-rc.d fan-daemon defaults >/dev/null
    invoke-rc.d fan-daemon start || exit $?
fi
